# overlays/pipewire-audio/home/ALL/.config/autostart/pipewire-pulse.desktop
[Desktop Entry]
Type=Application
Name=PipeWire pulse server autostart
Exec=pipewire-pulse
Terminal=False
