# overlay/pipewire-audio/home/ALL/.config/autostart/pipewire.desktop
[Desktop Entry]
Type=Application
Name=PipeWire server autostart
Exec=pipewire
Terminal=False
